# Guides/Raspberry Pi NodeRed & Monitor & Homekit2Mqtt/configs/homekit2mqtt/crontab
You can run homekit2mqtt on startup using cron. Do -
crontab -e
@reboot sudo homekit2mqtt --url mqtt://hass:XXXXXXXXXXXXX@192.168.79.54 --retain --name Homekit2Mqtt --username CC:22:3D:E3:CE:F1

